# node_network_details.socket
[Unit]
Description=Node Network Details Socket

[Socket]
ListenStream=9101
Accept=yes

[Install]
WantedBy=sockets.target
